# majordomo_nooliterx.conf
description	"MagorDoMo http://smartliving.ru/ Startup script for reciver NooLite"

start on runlevel [2345]
stop on runlevel [!2345]

#expect fork

script
  exec /usr/local/bin/nooliterx
end script
